# tautulli/tautulli.service
[Unit]
Description=Tautulli Service
After=plexmediaserver.service

[Service]
User=tautulli
Group=tautulli
Type=simple
ExecStart=/usr/bin/python /usr/lib/tautulli/PlexPy.py --nolaunch --config /var/lib/tautulli/config.ini --datadir /var/lib/tautulli
SyslogIdentifier=tautulli

[Install]
WantedBy=multi-user.target
