# python-fangfrisch/fangfrisch.timer
[Unit]
Description=Update unofficial ClamAV signatures
Documentation=https://rseichter.github.io/fangfrisch/
Requires=fangfrisch.service

[Timer]
# If you change the run frequency, make sure to match the settings
# in /etc/fangfrisch/fangfrisch.conf. See detailed comments there.
OnCalendar=*:0/10
RandomizedDelaySec=90
Unit=fangfrisch.service

[Install]
WantedBy=timers.target
