# ums/ums.timer
[Unit]
Description=Start UMS delayed on boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
