# mastodon/mastodon.target
[Unit]
Description=Mastodon
Wants=mastodon-web.service mastodon-sidekiq.service mastodon-streaming.service
Wants=postgresql.service valkey.service
Before=mastodon-web.service mastodon-sidekiq.service mastodon-streaming.service
After=postgresql.service valkey.service

[Install]
WantedBy=multi-user.target
