# backintime/backintime-backup-job.service
[Unit]
Description=Run backintime snapshot generation

[Service]
ExecStart=/usr/bin/backintime backup
