# yacy-git/yacy.service
[Unit]
Description=YaCy P2P Web Search
After=network.target

[Service]
Type=forking
ExecStart=/usr/share/java/yacy/startYACY.sh
ExecStop=/usr/share/java/yacy/stopYACY.sh
User=yacy
TimeoutStopSec=50

[Install]
WantedBy=multi-user.target
