# dahdi-linux-git/dahdi.service
[Unit]
Description=Configure DAHDI spans
Documentation=man:dahdi_cfg(8)
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service
ConditionPathExists=/etc/dahdi/system.conf

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/dahdi_cfg
ExecStop=/usr/bin/dahdi_cfg -s

[Install]
WantedBy=multi-user.target
