# octoprint-venv/octoprint.service
[Unit]
Description=3D Printing Web Server
After=network.target

[Service]
User=octoprint
Group=octoprint
Type=simple
ExecStart=/opt/octoprint-venv/bin/python -m octoprint serve
Restart=on-failure

[Install]
WantedBy=multi-user.target
