# sac-core/safenetauthenticationclient.service
[Unit]
Description=Thales eToken SACSrv daemon

[Service]
ExecStart=/usr/bin/SACSrv

[Install]
WantedBy=multi-user.target
