# lact/lactd.service
[Unit]
Description=AMDGPU Control Daemon
After=multi-user.target

[Service]
ExecStartPre=/usr/bin/sleep 2
ExecStart=/usr/bin/lact-daemon

[Install]
WantedBy=graphical.target
