# jackett/jackett.service
[Unit]
Description=Jackett Service
After=network.target

[Service]
User=jackett
Group=jackett
Restart=always
RestartSec=5
Type=simple
Environment=XDG_CONFIG_HOME=/var/lib/jackett
ExecStart=/usr/lib/jackett/jackett --NoRestart --NoUpdates --DataFolder /var/lib/jackett
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
