# rustdesk-server-bin/rustdesk-server-hbbr.service
[Unit]
Description=RustDesk Relay Server
Requires=network.target
After=systemd-user-sessions.service

[Service]
Type=simple
WorkingDirectory=/opt/rustdesk-server/data
ExecStart=/opt/rustdesk-server/hbbr
KillMode=mixed
TimeoutStopSec=30
User=rustdeskserver
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
