# hsphfpd-git/telephony_client.service
[Unit]
Description=telephony_client for hsphfpd

[Service]
Type=simple
ExecStart=/usr/bin/telephony_client.pl
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target

