# xampp/xampp.service.in
[Unit]
Description=XAMPP service (XAMPP Apache + MariaDB + PHP + Perl)

[Service]
Type=oneshot
ExecStart=@XAMPP_ROOT@/lampp start
ExecStop=@XAMPP_ROOT@/lampp stop
RemainAfterExit=yes
