# xbox-generic-controller/50_xbox-generic-controller.rules
ATTRS{idVendor}=="045e", ATTRS{idProduct}=="028e", RUN+="/usr/bin/fixcontroller.py"
