# trousers/tcsd.service
[Unit]
Description=Manager for Trusted Computing resources
After=sysinit.target

[Service]
User=tss
ExecStart=/usr/bin/tcsd -f -c /etc/tcsd.conf

[Install]
WantedBy=default.target
