# firefly-iii/firefly-iii-cron.timer
[Unit]
Description=Firefly III cron timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
