# xmm7360-pci-spat-dkms-git/xmm7360.service
[Unit]
Description=XMM7360 Modem Init
After=NetworkManager.service
Requires=multi-user.target systemd-user-sessions.service dev-ttyXMM2.device

[Service]
Type=oneshot
ExecStart=/usr/bin/lte up
RemainAfterExit=yes
TimeoutSec=120


[Install]
WantedBy=graphical.target
